# docker/php/Supervisor/larameus.conf
[program:sarvari]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/larameus/artisan horizon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=5
redirect_stderr=true
stdout_logfile=/var/www/larameus/storage/logs/larameus-sepervisor.log
stopwaitsecs=3600
